# systemd/awesome-release.path
[Unit]
Description="Monitor package build directory for new entries to release"

[Path]
PathModified=/srv/packages/deb
Unit=awesome-release.service

[Install]
WantedBy=multi-user.target
